# additional-modules/god/templates/redhat/init_d_prolog
#
# God
#
# chkconfig: - 85 15
# description: start, stop, restart God (bet you feel powerful)
#
